# plist/arthack.keeper-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for the keeper babysitter's EXTERNAL DEAD-MAN
  (keeper-watchdog) — "who watches the watcher".

  This is a TEMPLATE, not an installed agent: keeper ships no `install` verb.
  See README.md "Install" step 8b for the manual symlink + `launchctl bootstrap`.

  What it does: every 10 minutes, runs the standalone `keeper-watchdog`, which
  reads ONLY the babysitter's `heartbeat.json` (written as the last action on
  every completed keeper-watch tick) and ALARMS (notifyctl + botctl) when it has
  gone stale (> 900s). This is the ONE failure class the babysitter structurally
  cannot self-report: its own crash/hang. A dead monitor never runs its own
  monitor, so the watchdog is a SEPARATE launchd job — it does NOT open
  keeper.db, does NOT talk to keeperd, and does NOT depend on keeper-watch or
  keeperd being up (the exact case it exists to catch). Silent first-run (no
  heartbeat yet ≠ dead); once-daily all-clear so silence never means the
  watchdog itself died.

  Mirrors the orphanwatch/dropwatch dead-men: StartInterval 600, RunAtLoad,
  silent first-run, daily all-clear, PATH incl /Users/mike/.local/bin.

  Paths are hard-coded absolute because launchd does NOT source your shell rc,
  so `bun` cannot be found by name. The values below assume:
    - bun installed via Homebrew on Apple Silicon: /opt/homebrew/bin/bun
      (Intel Macs: /usr/local/bin/bun — edit ProgramArguments + PATH below)
    - the repo checked out at /Users/mike/code/keeper
    - babysitter heartbeat + watchdog logs under /Users/mike/.local/state/keeper-watch
  Edit these for your own username / arch before bootstrapping.

  PATH below MUST include /Users/mike/.local/bin: the alarm path shells out to
  `notifyctl` + `botctl`, which live there and would otherwise be unfound
  (launchd inherits no login PATH).

  Periodic via StartInterval (600s), NOT KeepAlive — mixing the two defeats the
  interval. RunAtLoad fires one immediate pass at bootstrap; the watchdog exits
  silently when there is no heartbeat yet, so an early load is harmless.

  Install:
    ln -s "$PWD/plist/arthack.keeper-watchdog.plist" ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/arthack.keeper-watchdog.plist

  Uninstall:
    launchctl bootout gui/$(id -u)/arthack.keeper-watchdog
    rm ~/Library/LaunchAgents/arthack.keeper-watchdog.plist

  Failures are observable via:
    launchctl print gui/$(id -u)/arthack.keeper-watchdog
    /Users/mike/.local/state/keeper-watch/watchdog.stderr
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>arthack.keeper-watchdog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/bun</string>
		<string>run</string>
		<string>/Users/mike/code/keeper/cli/keeper-watchdog.ts</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/mike/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>StartInterval</key>
	<integer>600</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>WorkingDirectory</key>
	<string>/Users/mike/code/keeper</string>
	<key>StandardOutPath</key>
	<string>/Users/mike/.local/state/keeper-watch/watchdog.stdout</string>
	<key>StandardErrorPath</key>
	<string>/Users/mike/.local/state/keeper-watch/watchdog.stderr</string>
</dict>
</plist>
