# plist/arthack.keeperd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for the keeper reducer daemon (keeperd).

  This is a TEMPLATE, not an installed agent: `scripts/install.sh` symlinks it
  into ~/Library/LaunchAgents/ and reloads it via `launchctl bootstrap`. See
  README.md "Install".

  Paths are hard-coded absolute because launchd does NOT source your shell rc,
  so `bun` cannot be found by name. The values below assume:
    - bun installed via Homebrew on Apple Silicon: /opt/homebrew/bin/bun
      (Intel Macs: /usr/local/bin/bun — edit ProgramArguments + PATH below)
    - the repo checked out at /Users/mike/code/keeper
    - state under /Users/mike/.local/state/keeper
  Edit these for your own username / arch before bootstrapping.

  NOTE: keeperd's transcript worker depends on the native `@parcel/watcher`
  addon, so the deploy environment MUST have `bun install` already run with a
  resolvable `node_modules/` under the WorkingDirectory (the repo checkout)
  before first bootstrap. The N-API addon load is a hard `dyld` crash if missing,
  not a catch-able error — keeperd will fail to boot, not degrade gracefully.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>arthack.keeperd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/bun</string>
		<string>run</string>
		<string>/Users/mike/code/keeper/src/daemon.ts</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<!-- ~/.local/bin carries the operator paging binaries (agentbot,
		  notifyctl); without it every daemon page spawn fails ENOENT and the
		  paging-channel-down probe latches. -->
		<string>/Users/mike/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<!--
		  Verbose server-worker diagnostic logging. Default `0` (off) writes
		  only the rare `[server-worker]` error class to server.stderr. Flip
		  to `1` and `launchctl kickstart -k gui/$UID/arthack.keeperd` to
		  enable `[srv-ts]` stage timings, frame byte counts, and connection
		  lifecycle on stderr. Read once at worker module load.
		-->
		<key>KEEPER_TRACE_SERVER</key>
		<string>0</string>
		<!-- Mandatory claude-swap CLI; independent of launchd's PATH. -->
		<key>KEEPER_CSWAP_BIN</key>
		<string>/Users/mike/.local/bin/cswap</string>
		<!--
		  The subscribe server's UDS defaults to keeperd.sock in the state dir
		  below (sibling of StandardOutPath / keeper.db). Uncomment to override:
		  <key>KEEPER_SOCK</key>
		  <string>/Users/mike/.local/state/keeper/keeperd.sock</string>
		-->
		<!--
		  The Agent Bus relay worker (epic fn-875) owns its OWN writable bus.db
		  and a DEDICATED relay socket bus.sock, both defaulting to the state dir
		  below (siblings of keeper.db / keeperd.sock, but a separate file +
		  protocol). Uncomment to override either:
		  <key>KEEPER_BUS_DB</key>
		  <string>/Users/mike/.local/state/keeper/bus.db</string>
		  <key>KEEPER_BUS_SOCK</key>
		  <string>/Users/mike/.local/state/keeper/bus.sock</string>
		-->
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
		<key>ThrottleInterval</key>
		<integer>10</integer>
	</dict>
	<!--
	  keeperd is the most latency-sensitive process on the host — it must fold
	  every hook event in real time — so it MUST NOT run in the throttled
	  Background QoS class, which starves it FIRST under host CPU contention.
	  Standard + Nice=-5 gives it scheduling priority without the all-throttling-
	  removed Interactive class (which can starve the human's foreground work).
	  launchd honors a negative Nice for its agents. ProcessType is read at SPAWN,
	  so after editing this re-read it via:
	    launchctl kickstart -k gui/$UID/arthack.keeperd
	  and confirm `ps -o pid,nice,comm -p <keeperd-pid>` shows nice -5.
	-->
	<key>ProcessType</key>
	<string>Standard</string>
	<key>Nice</key>
	<integer>-5</integer>
	<key>WorkingDirectory</key>
	<string>/Users/mike/code/keeper</string>
	<key>StandardOutPath</key>
	<string>/Users/mike/.local/state/keeper/server.stdout</string>
	<key>StandardErrorPath</key>
	<string>/Users/mike/.local/state/keeper/server.stderr</string>
</dict>
</plist>
